# scripts/sh/slackbot.cron
0 05 * * * ~/Downloads/slack/workspace/stopBots.sh
0 05 * * * ~/Downloads/slack/workspace/stopFAQBots.sh
0 05 * * * ~/Downloads/slack/taming-text/book/bin/stop-solr.sh

1 05 * * * ~/Downloads/slack/workspace/startBots.sh
1 05 * * * ~/Downloads/slack/workspace/startFAQBots.sh

1 05 * * * ~/Downloads/slack/taming-text/faqsetup/extract_setup.sh
2 05 * * * ~/Downloads/slack/taming-text/book/bin/start-faq-ws.sh

0 09 * * 5 python ~/Downloads/slack/workspace/send_mail.py